# templates/systemd.service.erb
# This file is managed by Puppet, any changes will be overwritten
[Unit]
Description=<%= @service_name %> web application.
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/usr/share/<%= @service_name %>
<% if @jvm_opts -%>
Environment="JAVA_OPTS=<%= @jvm_opts %>"
<% end -%>
ExecStart=/usr/share/<%= @service_name %>/bin/<%= @service_name %> -Dconfig.file=<%= @applicationconfig %> -Dpidfile.path=/dev/null -Dlogger.file=<%= @loggerconfig %>
Restart=on-failure
RestartSec=5
User=<%= @user %>
SuccessExitStatus=143
LimitNOFILE=65536
MemoryHigh=5632M
MemoryMax=6G
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
